# supervisor/download_server.ini
[program:download_server1]
command=/home/anaconda3/bin/python /usr/local/src/audio-download/start_app.py
directory=/usr/local/src/audio-download
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/home/log_sellbot/audio-download/audio_download.log