# supervisor/download_consumer.ini
[program:download_consumer1]
command=/home/anaconda3/bin/python /usr/local/src/audio-download/consumer/consumer.py
directory=/usr/local/src/audio-download/consumer
autostart=true
autorestart=true
user=root
stdout_logfile=none